# crontab
*/5 * * * * /bin/bash -l -c '/app/update.sh >> /app/cron.log 2>&1'